# buildquantumvis/pkg-fix/packages-templates/pipewire/REPLACEMENTS/etc/profile.d/pipewire
mkdir -p /run/pipewire
export PIPEWIRE_RUNTIME_DIR=/run/pipewire
